# .docker/php7/config/crontab
* * * * * cd /app && /usr/local/bin/php cron.php >>/app/log.txt 2>&1
